# backups/crontab
0 2 * * * /usr/local/bin/backup-web.sh >> /var/log/backup-web.log 2>&1
30 2 * * * /usr/local/bin/backup-database.sh >> /var/log/backup-db.log 2>&1